# util/xtendpid_daemon.service
[Unit]
Description=xtendpid daemon

[Service]
Type=forking
ExecStart=/usr/bin/xtendpid_daemon

[Install]
WantedBy=multi-user.target

