# util/xtendpid_daemon.service
[Unit]
Description=xtendpid daemon

[Service]
Type=forking
ExecStart=/usr/local/bin/xtendpid_daemon start

[Install]
WantedBy=multi-user.target

